# xdwim.service
[Unit]
Description=xdwim Service

[Service]
ExecStart=%h/.local/bin/xdwim
Restart=on-failure

[Install]
WantedBy=default.target